# templates/ai.openclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>

    <key>Comment</key>
    <string>OpenClaw Gateway with Auto-Recovery</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
      <key>NetworkState</key>
      <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProgramArguments</key>
    <array>
      <string>__NODE_PATH__</string>
      <string>__OPENCLAW_LIB_PATH__/node_modules/openclaw/dist/entry.js</string>
      <string>gateway</string>
      <string>--port</string>
      <string>__GATEWAY_PORT__</string>
    </array>

    <key>StandardOutPath</key>
    <string>__HOME_DIR__/.openclaw/logs/gateway.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/.openclaw/logs/gateway.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME_DIR__</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
